# crontab
# m h  dom mon dow   command
* * * * * sleep 30; ~/update-soundrepo.sh
* * * * * ~/update-soundrepo.sh
